# scripts/com.studyg.generate-pages.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.studyg.generate-pages</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/michalrucka/Desktop/Projects/Study G/StudyG-blog/scripts/generate-pages.sh</string>
  </array>

  <!--
    3 pages/day: 09:00, 14:00, 20:00
    launchd fires on the NEXT scheduled time after a missed window (unlike cron).
    The script's 6-hour guard prevents double-runs if the Mac wakes late.
  -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>14</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>20</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/michalrucka/Desktop/Projects/Study G/StudyG-blog/scripts/generate-pages.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/michalrucka/Desktop/Projects/Study G/StudyG-blog/scripts/generate-pages.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>

</dict>
</plist>
